# supervisord.conf
[supervisord]
nodaemon=true

[program:docker_startup]
command=/home/zadatak/mojizadaci/docker_startup.sh
autostart=true
autorestart=false
startretries=0
redirect_stderr=false
redirect_stderr=false


